# Linux/Debian/jessie/system_root/etc/profile.d/proxy.sh
export http_proxy=http://172.16.81.1:3128/
export https_proxy=http://172.16.81.1:3128/
